# build/s6/softether-vpnsrv/run
#!/bin/bash

# SOFTETHER START
echo "Starting Softether Server..."
echo "-------------------------------"
/s6-bin/softether-vpnsrv/vpnserver start

# Create IP tables rules
echo "Creating postrouting rules."
iptables -t nat -A POSTROUTING -s ${SRVIPSUBNET:-10.0.0}.0/24 -j MASQUERADE

# Create DHCP configuration
if [ -f /cfg/dnsmasq.conf ]; then
  echo "Creating DHCP server configuration..."
  sed "s/\$SRVIPSUBNET/${SRVIPSUBNET:-10.0.0}/g" /cfg/dnsmasq.conf >/etc/dnsmasq.conf
fi

# Link config files
if [ -f /cfg/vpn_server.config ]; then
  echo "Linking configuration file."
  rm /etc/vpnserver/vpn_server.config 2>/dev/null
  ln -s /cfg/vpn_server.config /etc/vpnserver
fi

# Start DHCP Server
echo "Starting DHCP server."
/usr/sbin/dnsmasq

# Bind to created tap interface
echo "Waiting 3 seconds."
sleep 3
echo "Binded IP address to the server with ${SRVIPSUBNET:-10.0.0}.1"
/sbin/ifconfig tap_soft ${SRVIPSUBNET:-10.0.0}.1

# Check health
while ping -c 1 -W 10 "${SRVIPSUBNET:-10.0.0}.1" >&/dev/null && pgrep dnsmasq >&/dev/null; do
  echo "-------------------------------"
  echo "Self test success. VPN Server alive and can access ${SRVIPSUBNET:-10.0.0}.1."
  if [[ -z "${KEEP_SERVER_LOG}" ]]; then
    echo "Clearing up server log files."
    rm /usr/local/libexec/softether/vpnserver/server_log/*.log -f
  fi
  if [[ -z "${KEEP_PACKET_LOG}" ]]; then
    echo "Clearing up packet log files."
    rm /usr/local/libexec/softether/vpnserver/packet_log/*.log -f
  fi
  if [[ -z "${KEEP_SECURITY_LOG}" ]]; then
    echo "Clearing up security log files."
    rm /usr/local/libexec/softether/vpnserver/packet_log/*.log -f
  fi
  echo "Sleeping ${SLEEPTIME:-600} seconds."
  sleep ${SLEEPTIME:-600}
done
