# build/services.d/softether-vpnsrv/finish
#!/bin/bash

echo "-------------------------------"
echo "Terminating softether-vpnsrv."
echo "-------------------------------"
echo "Shutting down the daemon"
/s6-bin/softether-vpnsrv/vpnserver stop >&/dev/null
echo "Stoping tap adapter."
ip tuntap del dev tap_soft mode tap
ifconfig tap_soft down >&/dev/null
echo "Deleting Virtual Network Adapter"
ip link delete tap_soft >&/dev/null
echo "Shutting down DHCP server."
pkill /usr/sbin/dnsmasq >&/dev/null
echo "-------------------------------"
echo "Terminated softether-vpnsrv."
echo "-------------------------------"
