# build/services.d/softether-vpnsrv/finish
#!/bin/bash

source /scripts/logger.sh

log_start "Terminating services..." "both"

log_warn "Shutting down SoftEther daemon..."
/s6-bin/softether-vpnsrv/vpnserver stop >&/dev/null

log_warn "Removing tap software adapter..."
ip tuntap del dev tap_soft mode tap
ifconfig tap_soft down >&/dev/null

log_warn "Deleting SoftEther virtual network adapter..."
ip link delete tap_soft >&/dev/null

log_warn "Shutting down DNSMASQ daemon..."
pkill /usr/sbin/dnsmasq >&/dev/null

log_finish "Terminated services." "both"
