# build/services.d/softether-vpnsrv/run
#!/bin/bash

# SOFTETHER START
echo "Starting Softether Server..."
echo "-------------------------------"
/s6-bin/softether-vpnsrv/vpnserver start

# Start DHCP Server
echo "Starting DHCP server."
/usr/sbin/dnsmasq

# Bind to created tap interface
echo "Waiting 3 seconds."
s6-sleep 3
echo "Binded IP address to the server with ${SRVIPSUBNET:-10.0.0}.1"
ip tuntap add dev tap_soft mode tap
/sbin/ifconfig tap_soft ${SRVIPSUBNET:-10.0.0}.1

# Check health
while ping -c 1 -W 10 "${SRVIPSUBNET:-10.0.0}.1" >&/dev/null && pgrep dnsmasq >&/dev/null; do
  echo "-------------------------------"
  echo "Self test success. VPN Server alive and can access ${SRVIPSUBNET:-10.0.0}.1."
  if [[ -z "${KEEP_SERVER_LOG}" ]]; then
    echo "Clearing up server log files."
    rm /usr/local/libexec/softether/vpnserver/server_log/*.log -f
  fi
  if [[ -z "${KEEP_PACKET_LOG}" ]]; then
    echo "Clearing up packet log files."
    rm /usr/local/libexec/softether/vpnserver/packet_log/*.log -f
  fi
  if [[ -z "${KEEP_SECURITY_LOG}" ]]; then
    echo "Clearing up security log files."
    rm /usr/local/libexec/softether/vpnserver/packet_log/*.log -f
  fi
  echo "Sleeping ${SLEEPTIME:-600} seconds."
  s6-sleep ${SLEEPTIME:-600}
done
