# build/services.d/softether-vpnsrv/run
#!/bin/bash

source /scripts/logger.sh

# SOFTETHER START
log_start "Starting Softether Server..."
/s6-bin/softether-vpnsrv/vpnserver start

# Start DHCP Server
log_start "Starting DHCP server."
/usr/sbin/dnsmasq

log_divider

# Bind to created tap interface
log_wait "Waiting 3 seconds."
s6-sleep 3

# Bind to ip address with tap device
log_info "Binded IP address to the server with ${SRVIPSUBNET:-10.0.0}.1"
ip tuntap add dev tap_soft mode tap
/sbin/ifconfig tap_soft ${SRVIPSUBNET:-10.0.0}.1 netmask ${SRVIPNETMASK:-255.255.255.0}

# Check health
while ping -c 1 -W 10 "${SRVIPSUBNET:-10.0.0}.1" >&/dev/null && pgrep dnsmasq >&/dev/null; do
  log_info "[$(date +"%Y%m%d-%H:%M:%S")] Self test success. VPN Server alive and can access ${SRVIPSUBNET:-10.0.0}.1." "top"

  if [[ -z "${KEEP_SERVER_LOG}" ]]; then
    log_info "Clearing up server log files."
    rm /usr/local/libexec/softether/vpnserver/server_log/*.log -f
  fi

  if [[ -z "${KEEP_PACKET_LOG}" ]]; then
    log_info "Clearing up packet log files."
    rm /usr/local/libexec/softether/vpnserver/packet_log/*.log -f
  fi

  if [[ -z "${KEEP_SECURITY_LOG}" ]]; then
    log_info "Clearing up security log files."
    rm /usr/local/libexec/softether/vpnserver/packet_log/*.log -f
  fi

  log_wait "Next check in ${SLEEPTIME:-600} seconds."
  s6-sleep ${SLEEPTIME:-600}
done
